# resources/puppet/site/websockify/templates/websockify.service.epp
[Unit]
Description=websockify.service
After=network.target

[Service]
ExecStart=<%= $websockify::_cmd %>
Restart=on-failure

[Install]
WantedBy=default.target
